# deploy/common/files/etc/systemd/system/babytracker.service
[Unit]
Description=BabyTracker Web Application
After=network.target postgresql.service
Requires=postgresql.service
ConditionPathExists=!/var/lib/babytracker/.needs-setup

[Service]
Type=simple
User=babytracker
Group=babytracker
EnvironmentFile=/etc/babytracker/babytracker.env
ExecStart=/usr/local/bin/babytracker
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Allow binding to privileged ports (HTTPS on :443 and HTTP-01 on :80).
# Note: we intentionally do NOT set CapabilityBoundingSet or NoNewPrivileges —
# the self-update handler invokes `sudo systemctl restart babytracker`, and
# sudo needs CAP_SETUID/CAP_SETGID that a restrictive bounding set would
# block.
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Security hardening. /usr/local/bin is writable so self-update can atomically
# rename the new binary over the old one; the binary file is owned by the
# babytracker user (see layer/babytracker.yaml).
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/babytracker /usr/local/bin
PrivateTmp=true

[Install]
WantedBy=multi-user.target
