# deploy/rpi-image-gen/files/etc/systemd/system/babytracker-setup-ap.service
[Unit]
Description=BabyTracker Setup Access Point and Portal
After=babytracker-firstboot.service network.target
Requires=babytracker-firstboot.service
ConditionPathExists=/var/lib/babytracker/.needs-setup

[Service]
Type=simple

# Start AP infrastructure before the app.
# `+` prefix → run as root (overrides User= for that command).
# `-` prefix → ignore non-zero exit (best-effort).
ExecStartPre=+-/usr/bin/ip addr add 192.168.4.1/24 dev wlan0
ExecStartPre=+/usr/bin/systemctl start hostapd
ExecStartPre=+/usr/bin/systemctl start dnsmasq
ExecStartPre=+/usr/sbin/iptables -t nat -A PREROUTING -i wlan0 -p tcp --dport 80 -j REDIRECT --to-port 8099

# Run the BabyTracker binary in setup mode (serves setup wizard)
User=babytracker
Group=babytracker
EnvironmentFile=/etc/babytracker/babytracker.env
ExecStart=/usr/local/bin/babytracker

# Allow binding to privileged ports (the binary listens on :80 in setup mode)
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Cleanup when stopped — also need root for these
ExecStopPost=+-/usr/sbin/iptables -t nat -D PREROUTING -i wlan0 -p tcp --dport 80 -j REDIRECT --to-port 8099
ExecStopPost=+-/usr/bin/systemctl stop hostapd
ExecStopPost=+-/usr/bin/systemctl stop dnsmasq
ExecStopPost=+-/usr/bin/ip addr del 192.168.4.1/24 dev wlan0

Restart=no

[Install]
WantedBy=multi-user.target
