# deploy/rpi-image-gen/files/etc/systemd/system/babytracker-setup-ap.service
[Unit]
Description=BabyTracker Setup Access Point and Portal
After=babytracker-firstboot.service network.target
Requires=babytracker-firstboot.service
ConditionPathExists=/var/lib/babytracker/.needs-setup

[Service]
Type=simple

# Conditionally start the captive portal AP — skipped if the device already
# has a network connection (e.g. Ethernet with DHCP). `+` prefix runs as root.
ExecStartPre=+/usr/local/bin/babytracker-maybe-start-ap.sh

# Run the BabyTracker binary in setup mode (serves setup wizard on :443 + :80)
User=babytracker
Group=babytracker
EnvironmentFile=/etc/babytracker/babytracker.env
ExecStart=/usr/local/bin/babytracker

# Allow binding to privileged ports (:443 and :80).
# CapabilityBoundingSet is intentionally NOT set: setup mode invokes sudo to
# call setup-wifi.sh / setup-ethernet.sh, and sudo needs CAP_SETUID/CAP_SETGID
# which a restrictive bounding set would block.
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Cleanup when stopped (only tears down what was brought up)
ExecStopPost=+/usr/local/bin/babytracker-maybe-stop-ap.sh

Restart=no

[Install]
WantedBy=multi-user.target
